# etc/crontab
0 5  * * *  /root/elasticsearch/cron_job.sh
